# sparse/usr/lib/systemd/system/sailfish-unlock-agent.service.d/50-vendor.hwcomposer-2-1.conf
[Service]
# stop hwcomposer before unlock ui
ExecStartPre=-/system/bin/stop vendor.hwcomposer-2-1

# start hwcomposer after unlock ui
ExecStart=-/system/bin/start vendor.hwcomposer-2-1
