# sparse/usr/lib/systemd/system/jolla-actdead-charging.service.d/10-wait_for_hwcomposer.conf
[Service]
ExecStartPre=/usr/bin/droid/wait_for_hwcomposer
